# deploy/systemd/infraguardian-collector.service
[Unit]
Description=InfraGuardian Host Collector
After=network.target

[Service]
Type=simple
User=igadmin
WorkingDirectory=/opt/infraguardian
EnvironmentFile=/opt/infraguardian/collector/config.env
ExecStart=/opt/infraguardian/venv/bin/python -u /opt/infraguardian/collector/collector.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
